# runner/app/runit/infer/run
#!/bin/sh

echo "Starting infer"
exec 2>&1
exec python -u /app/infer/infer.py